# ansible/roles/openclaw/templates/openclaw.service.j2
[Unit]
Description=OpenClaw AI Agent Gateway - {{ agent_name }}
After=network.target ollama.service
Requires=ollama.service

[Service]
Type=simple
User={{ openclaw_service_user }}
WorkingDirectory={{ openclaw_install_dir }}
Environment="NODE_ENV=production"
EnvironmentFile={{ openclaw_install_dir }}/.env
ExecStart=/usr/bin/openclaw gateway --config {{ openclaw_install_dir }}/openclaw.json --port {{ openclaw_dashboard_port }}
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
