# ansible/roles/proxmox_fan_control/templates/dell-fan-control.service.j2
[Unit]
Description=Dell PowerEdge Fan Speed Control
Documentation=https://github.com/infiquetra/home-lab
# Run after the IPMI service is available
After=ipmi.service systemd-modules-load.service
# Not a hard dependency — the script handles failures gracefully
Wants=ipmi.service

[Service]
Type=oneshot
# Re-apply fan settings on every start (handles reboots)
ExecStart={{ fan_control_script }}
RemainAfterExit=yes
# Do not restart on failure — iDRAC 9 firmware 7.x blocks these commands
# (insufficient privilege level is a firmware restriction, not a transient error)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
